# other/usr/lib/systemd/system/reflector.timer
[Unit]
Description=Daily pacman mirrorlist refresh

[Timer]
OnCalendar=Daily
Persistent=true
AccuracySec=5m

[Install]
WantedBy=timers.target
